# conf/systemd.service
[Unit]
Description=__APP__ daemon
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
Environment=NODE_ENV=production
Environment=NODE_CONFIG_DIR=__FINALPATH__/config
User=__APP__
Group=__APP__
ExecStart=/bin/sh -c '/usr/bin/npm start'
WorkingDirectory=__FINALPATH__/
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=__APP__
Restart=always

[Install]
WantedBy=multi-user.target
